# battmon.service
[Unit]
Description=Battery monitor with notifications for my ASUS laptop

[Service]
Type=oneshot
ExecStart=/usr/local/bin/battmon
